# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton-rootfs/usr/lib/systemd/system/afterburn-checkin.service.d/10-override.conf
# We don't set the cloud provider on the kernel command line. The unit will only
# run on azure
[Unit]
ConditionKernelCommandLine=
ConditionVirtualization=microsoft

[Service]
ExecStart=
ExecStart=/usr/bin/afterburn --provider=azure --check-in
