# libvirt/image/files/etc/systemd/system/run-kata\x2dcontainers-shared-containers.mount
[Unit]
Description=Mount unit for /run/kata-containers/shared/containers
Before=kata-agent.service

[Mount]
What=tmpfs
Where=/run/kata-containers/shared/containers
Type=tmpfs
Options=mode=755

[Install]
WantedBy=multi-user.target
