# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/process-user-data.service.d/10-override.conf
[Service]
# mount config disk if available
ExecStartPre=-/bin/mount -t iso9660 -o ro /dev/disk/by-label/cidata /media/cidata
# Extend initdata digest into TDX RTMR 3.
# Fails the service if the TDX sysfs node is missing or the write fails —
# running without RTMR3 measurement would break the attestation chain.
#
# 1. Read the SHA-384 hex digest (64–96 hex chars).
# 2. Right-pad to 96 hex chars (48 bytes) with zeros.
# 3. Convert hex pairs to \xHH escape sequences via sed, then printf
#    interprets them into raw bytes.
# 4. dd bs=48 count=1 guarantees a single 48-byte write(2) syscall,
#    which the TDX kernel interface requires.
#
# Note: $$ and %% are systemd escaping (literal $ and %).
#       \\\\x becomes \\x after systemd+shell, so sed emits \xHH.
ExecStartPost=/bin/sh -c '\
  hex=$(cat /run/peerpod/initdata.digest); \
  hex=$(printf "%%-96s" "$$hex" | tr " " "0"); \
  printf "$$(echo "$$hex" | sed "s/../\\\\x&/g")" | \
  dd of=/sys/devices/virtual/misc/tdx_guest/measurements/rtmr3:sha384 bs=48 count=1 2>/dev/null'
