# src/cloud-api-adaptor/podvm/mkosi.skeleton-sftp/usr/lib/systemd/system/ssh.service.d/10-host-keys.conf
[Unit]
Wants=ssh-host-keys.service
After=ssh-host-keys.service
