# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/afterburn-checkin.service.d/10-override.conf
# As our image is generic, we don't set the cloud provider on the kernel command line.
# Instead, we always run the unit, even if it is only needed on Azure right now.
[Unit]
ConditionKernelCommandLine=

[Service]
ExecStart=
ExecStart=-/usr/bin/afterburn --provider=azure --check-in
