# src/cloud-api-adaptor/podvm/files/etc/systemd/system/agent-protocol-forwarder.path
[Unit]
Description=Wait for kata-agent socket before starting Agent Protocol Forwarder

[Path]
PathExists=/run/kata-containers/agent.sock
Unit=agent-protocol-forwarder.service

[Install]
WantedBy=multi-user.target
