# src/cloud-api-adaptor/podvm/files/etc/systemd/system/attestation-agent.path
[Unit]
Description=Watch for attestation agent configuration file

[Path]
PathExists=/run/peerpod/aa.toml
Unit=attestation-agent.service

[Install]
WantedBy=multi-user.target