# src/cloud-api-adaptor/podvm/files/etc/systemd/system/confidential-data-hub.path
[Unit]
Description=Monitor for the Attestation Agent socket

[Path]
PathExists=/run/confidential-containers/attestation-agent/attestation-agent.sock
Unit=confidential-data-hub.service

[Install]
WantedBy=multi-user.target
