# src/cloud-api-adaptor/podvm/files/etc/systemd/system/netns@.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=ip netns add %I
ExecStart=ip netns exec %I ip link set lo up
ExecStop=ip netns del %I
