# src/cloud-api-adaptor/podvm/files/etc/systemd/system/nvidia-cdi.service
[Unit]
Description=Generate NVIDIA CDI Configuration
# Must run after DKMS builds and loads the NVIDIA kernel modules.
After=nvidia-dkms-build.service nvidia-persistenced.service
Requires=nvidia-dkms-build.service nvidia-persistenced.service
Before=kata-agent.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate-nvidia-cdi.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
