# src/cloud-api-adaptor/podvm/files/etc/systemd/system/process-user-data.path
[Unit]
Description=Monitor for the user data file path

[Path]
PathExists=/media/cidata/user-data
Unit=process-user-data.service

[Install]
WantedBy=multi-user.target
