# src/cloud-api-adaptor/podvm/files/etc/systemd/system/scratch-storage.path
[Unit]
Description=Watch for scratch storage marker

[Path]
PathExists=/run/peerpod/scratch-space.marker
Unit=scratch-storage.service

[Install]
WantedBy=multi-user.target
