# src/cloud-api-adaptor/podvm/files/etc/systemd/system/nvidia-persistenced.service.d/10-after-modules.conf
[Unit]
After=systemd-modules-load.service
