# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/etc/profile.d/10-alias.sh
#!/usr/bin/env bash

alias neofetch='neofetch \
      --config /etc/neofetch/config.conf \
      --ascii /etc/neofetch/coco.ascii \
      --colors 9 9 9 9 9 12 \
      --ascii_colors 24 9'
alias pcrs='tpm2_pcrread sha256:all'
