# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/etc/profile.d/10-alias.sh
#!/usr/bin/env bash

alias fastfetch='fastfetch \
      --config /etc/fastfetch/config.jsonc'
alias pcrs='tpm2_pcrread sha256:all'
