# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/nvidia-persistenced.service.d/10-retry.conf
# nvidia-persistenced requires the nvidia kernel module to be loaded.
# On first boot, DKMS compiles the module from source (dkms.service), and
# systemd-modules-load.service reads /usr/lib/modules-load.d/nvidia.conf to
# insert the compiled modules into the running kernel. nvidia-persistenced must
# wait for both: dkms to finish building and modules-load to finish inserting.
# nvidia-cdi.service depends on nvidia-persistenced, so it also fails if this
# ordering is wrong — meaning no CDI spec is generated and GPUs are invisible.
[Unit]
After=dkms.service systemd-modules-load.service

[Service]
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=10
