# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/nvidia-persistenced.service.d/10-retry.conf
# nvidia-persistenced requires NVIDIA kernel modules to be loaded.
# On first boot DKMS compiles the modules (dkms.service). systemd-modules-load
# runs too early (before DKMS), so we load the modules explicitly with
# ExecStartPre after DKMS has finished building them.
[Unit]
After=dkms.service
StartLimitIntervalSec=60s
StartLimitBurst=10

[Service]
ExecStartPre=/sbin/modprobe nvidia
ExecStartPre=/sbin/modprobe nvidia-uvm
ExecStartPre=/sbin/modprobe nvidia-modeset
Restart=on-failure
RestartSec=5s
