# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/nvidia-persistenced.service.d/10-retry.conf
# nvidia-persistenced requires the nvidia kernel module to be loaded.
# On first boot, DKMS compiles the module from source, so nvidia-persistenced
# must wait for dkms.service to complete before starting. Without this ordering
# constraint both services start in parallel at basic.target and
# nvidia-persistenced fails because nvidia.ko isn't loaded yet.
# nvidia-cdi.service depends on nvidia-persistenced, so it also fails,
# meaning no CDI spec is generated and GPUs are not visible in containers.
[Unit]
After=dkms.service

[Service]
Restart=on-failure
RestartSec=3s
StartLimitIntervalSec=60s
StartLimitBurst=10
