# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton-sftp/usr/lib/systemd/system/sshd-keygen@.service.d/10-override.conf
# generate server SSH keys in /run
[Unit]
ConditionFileNotEmpty=
ConditionFileNotEmpty=|!/run/ssh/ssh_host_%i_key

[Service]
ExecStartPre=mkdir -p /run/ssh
ExecStart=
ExecStart=/usr/libexec/openssh/sshd-keygen %i -f /run/ssh/ssh_host_%i_key
