# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/nvidia-cdi-refresh.service.d/override.conf
[Unit]
After=nvidia-persistenced.service nvidia-fabricmanager.service
Wants=nvidia-persistenced.service nvidia-fabricmanager.service

[Service]
ExecCondition=/usr/local/bin/check-nvidia-gpu
