# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/agent-protocol-forwarder.service.d/10-wait-for-user-data.conf
# Ensure agent-protocol-forwarder starts only after process-user-data has
# written /etc/default/agent-protocol-forwarder with the OPTIONS env var
# (which contains the CAA proxy address). Without this ordering, APF starts
# before user-data is processed and fails with empty OPTIONS.
[Unit]
After=process-user-data.service
