# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/nvidia-fabricmanager.service.d/override.conf
[Unit]

[Service]
ExecCondition=/usr/local/bin/check-nvidia-gpu
RuntimeDirectory=nvidia-fabricmanager
Environment=FM_CONFIG_FILE=/run/nvidia-fabricmanager/fabricmanager.cfg
ExecStartPre=/bin/bash -c 'sed "s/^LOG_USE_SYSLOG=0/LOG_USE_SYSLOG=1/" /usr/share/nvidia/nvswitch/fabricmanager.cfg > /run/nvidia-fabricmanager/fabricmanager.cfg'
