# podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/process-user-data.service.d/10-override.conf
# On a read-only fs the kata-agent config is created in /run/peerpod, since it contains
# a parameter that can be set at pod creation time.
[Service]
ExecStart=
ExecStart=/usr/local/bin/process-user-data update-agent-config --agent-config-file /run/peerpod/agent-config.toml
