# src/cloud-api-adaptor/podvm-mkosi/mkosi.skeleton-rootfs/etc/profile.d/20-ssh-banner.sh
#!/usr/bin/env bash

if [[ -n ${SSH_CONNECTION} && $- == *i* ]] ; then
    neofetch
fi
