# podvm-mkosi/mkosi.skeleton/usr/lib/systemd/system/kata-agent.service.d/10-override.conf
# On a read-only fs the kata-agent config is created in /run/peerpod, since it contains
# a parameter that can be set at pod creation time.
[Service]
Environment=KATA_AGENT_CONFIG_PATH=/run/peerpod/agent-config.toml
ExecStart=
ExecStart=/usr/local/bin/kata-agent --config ${KATA_AGENT_CONFIG_PATH}
ExecStop=
ExecStopPost=/usr/local/bin/kata-agent-clean --config ${KATA_AGENT_CONFIG_PATH}
